# dropbox_sync_aperops.ini
[group:aperops_tools]
programs=dropbox_sync

[program:dropbox_sync]
command=/bin/bash -c 'cd /home/apertif/operations/aperops; ./dropbox_sync.py'
environment=HOME=/home/apertif
user=apertif
stopsignal=KILL
stopasgroup=true ; also kills the subprocesses
stdout_logfile=/data/apertif/var/log/dropbox_sync.log
redirect_stderr=true
stderr_logfile=NONE
